# deployments/systemd/tz.service
[Unit]
Description=TZ - Lightweight Linux Monitoring Agent
Documentation=https://github.com/Alice/TZ
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=tz
Group=tz
ExecStart=/usr/local/bin/tz serve --config /etc/tz/tz.yaml

# Restart policy
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log/tz
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes

# Resource limits
LimitNOFILE=65536
MemoryHigh=35M
MemoryMax=50M
CPUQuota=5%
TasksMax=50

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tz

[Install]
WantedBy=multi-user.target
